# Koha/Plugin/Com/Theke/SLNP/scripts/slnp-server.service
[Unit]
Description=SLNP server.

[Service]
ExecStart=/usr/sbin/slnp-server.sh $(koha-list)

[Install]
WantedBy=multi-user.target